# data/io.github.up-check.service.in
[Unit]
Description=Up — background update check
Documentation=https://github.com/VictoryTek/Up

[Service]
Type=oneshot
ExecStart=@BINDIR@/up --check
